# root-config/systemd-boot-update.service.d/override.conf
[Service]
ExecStartPost=/usr/bin/sbctl sign-all
