# root-config/checkupdates.timer
[Unit]
Description=Download package updates hourly

[Timer]
OnBootSec=1min
OnUnitInactiveSec=1hour

[Install]
WantedBy=timers.target
